# hivewrightv2-dashboard.service
[Unit]
Description=HiveWright v2 Dashboard (Next.js)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/hivewrightv2
ExecStart=/usr/bin/npm run dev
Restart=always
RestartSec=10
Environment=HIVEWRIGHT_RUNTIME_ROOT=%h/.hivewright
EnvironmentFile=%h/.hivewright/config/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
