# packaging/systemd/hivewright-dispatcher.service.example
[Unit]
Description=HiveWright Dispatcher
After=network.target hivewright-dashboard.service

[Service]
Type=simple
WorkingDirectory=%h/hivewright
EnvironmentFile=%h/hivewright/.env
ExecStart=/bin/bash %h/hivewright/start-dispatcher.sh
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
